# infra/systemd/spark-a/spark-a.service
[Unit]
Description=Spark A — Qwen3.6-35B-A3B-FP8 (LUMEOS Orchestrator)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/home/admin/start-spark-a.sh
ExecStop=/usr/bin/docker stop vllm-qwen
Restart=always
RestartSec=5
TimeoutStartSec=300
TimeoutStopSec=30

# Allow plenty of file descriptors for vLLM internal pools
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
